# contrib/hfast.socket
[Unit]
Description=HTTPd Fast Socket

[Socket]
ListenStream=80
ListenStream=443
ListenDatagram=443

[Install]
WantedBy=sockets.target
